# act/tasksched/zelin-webui.xml
<?xml version="1.0" encoding="UTF-8"?>
<!-- TEMPLATE: rendered by install.ps1 (see act/lib/taskscheduler.py), registered
     as \ZelinAIAssistant\webui.

     webui — the local web dashboard (act/webui.py, shipped v0.24.0). On Windows
     this IS the UI: the Mac menu-bar app cannot run here, so the browser board is
     how a Windows user reads state/dashboard.json and approves inbox cards.
     Resident like actd so it auto-starts at logon and stays up (LogonTrigger +
     RestartOnFailure). The Windows mirror of act/systemd/zelin-webui.service.

     SECURITY: act.webui binds 127.0.0.1 ONLY and requires a per-install token
     (state/webui.token, 0600) on every /api/* request — see act/webui.py. This
     task adds nothing network-facing.

     URL: act.webui prints "web dashboard on http://127.0.0.1:<port>" to its
     stdout; open that in a browser on the same box. -->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Zelin AI Assistant — local web dashboard (the Windows UI)</Description>
    <URI>\ZelinAIAssistant\webui</URI>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
    <RestartOnFailure>
      <Interval>PT1M</Interval>
      <Count>999</Count>
    </RestartOnFailure>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -NonInteractive -Command "$env:AIASSISTANT_HOME='@REPO_ROOT@'; $env:PATH='@CLAUDE_BIN_DIR@;'+$env:PATH; &amp; '@PYTHON@' -m act.webui"</Arguments>
      <WorkingDirectory>@REPO_ROOT@</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
