# act/tasksched/zelin-weekly-digest.xml
<?xml version="1.0" encoding="UTF-8"?>
<!-- TEMPLATE: rendered by install.ps1 (see act/lib/taskscheduler.py), registered
     as \ZelinAIAssistant\weekly-digest. Repetition every hour (mirror of the
     systemd zelin-weekly-digest.timer hourly wake).

     weekly digest (CONTRACT §24) — files a review-lane digest card plus a few
     automation-proposal cards from the last 7 days of ingest output. The task
     wakes it hourly; the module itself gates on sources.weekly_digest
     (enabled/day/hour, default Monday 09:xx) + a run-once-per-week marker, so
     schedule changes in config take effect WITHOUT re-registering the task. -->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Zelin AI Assistant — weekly digest (hourly wake, module gates to Mon 09:xx)</Description>
    <URI>\ZelinAIAssistant\weekly-digest</URI>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
      <Repetition>
        <Interval>PT1H</Interval>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT1H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -NonInteractive -Command "$env:AIASSISTANT_HOME='@REPO_ROOT@'; $env:PATH='@CLAUDE_BIN_DIR@;'+$env:PATH; &amp; '@PYTHON@' -m act.weekly_digest"</Arguments>
      <WorkingDirectory>@REPO_ROOT@</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
