# act/launchd/com.zelin.aiassistant.slackradar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TEMPLATE: install.sh renders the /Users/YOURUSERNAME placeholders (python,
  repo root, $HOME) into ~/Library/LaunchAgents — do not load this file as-is.

  slackradar — captures Slack DMs / group DMs / watched-channel @mentions and
  turns the actionable ones into requirement cards (draft-only for comms).

  SAFE UNDER LAUNCHD: this radar only touches the network + ~/Projects/zelin-ai-assistant/state,
  never ~/Documents, so it is NOT TCC-blocked (unlike the Obsidian radar, which
  must use crontab). It no-ops until a Slack user token exists at the path in
  config (sources.slack_token_path, default ~/Desktop/Keys/slack-user-token.txt).

  TOKEN: needs a Slack **user token** (xoxp-) with scopes
  search:read, im:history, im:read, mpim:history, mpim:read, channels:history,
  groups:history, users:read. See docs/SLACK_SETUP.md.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zelin.aiassistant.slackradar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOURUSERNAME/miniconda3/bin/python3</string>
        <string>-m</string>
        <string>act.radar_slack</string>
        <string>--once</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOURUSERNAME/Projects/zelin-ai-assistant</string>

    <!-- poll every 3 minutes -->
    <key>StartInterval</key>
    <integer>180</integer>
    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AIASSISTANT_HOME</key>
        <string>/Users/YOURUSERNAME/Projects/zelin-ai-assistant</string>
        <key>PATH</key>
        <string>/Users/YOURUSERNAME/.claude-bin:/opt/homebrew/bin:/opt/homebrew/sbin:/Users/YOURUSERNAME/miniconda3/bin:/Users/YOURUSERNAME/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOURUSERNAME/Projects/zelin-ai-assistant/state/slackradar.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOURUSERNAME/Projects/zelin-ai-assistant/state/slackradar.launchd.log</string>
</dict>
</plist>
