# act/launchd/com.zelin.aiassistant.radar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TEMPLATE: install.sh renders the /Users/YOURUSERNAME placeholders (python,
  repo root, $HOME) into ~/Library/LaunchAgents — do not load this file as-is.

  radar — periodic requirement scanner (every 30 min) for Zelin's AI Assistant.
  Scans incremental sources (Obsidian "2 - raw", watched Slack) and updates the
  requirement registry (act/registry/*.yaml), emitting new approval cards.

  ============================ RECOMMENDED: use crontab ============================
  launchd is TCC-BLOCKED from ~/Documents, and radar MUST read
  "~/Documents/Obsidian Vault/2 - raw". A crontab job inherits Full Disk Access
  (once Terminal/cron is granted it in System Settings > Privacy), so PREFER cron.

  Example crontab line (run `crontab -e` and paste; every 30 minutes):

    */30 * * * * cd /Users/YOURUSERNAME/Projects/zelin-ai-assistant && AIASSISTANT_HOME=/Users/YOURUSERNAME/Projects/zelin-ai-assistant PATH=/Users/YOURUSERNAME/.local/bin:/opt/homebrew/bin:/Users/YOURUSERNAME/miniconda3/bin:/usr/local/bin:/usr/bin:/bin /usr/bin/env python3 -m act.radar >> /Users/YOURUSERNAME/Projects/zelin-ai-assistant/state/radar.cron.log 2>&1

  (~/.local/bin leads that PATH on purpose: it is where the claude CLI usually
  lives — an outdated second install in /opt/homebrew/bin ranked first once
  broke dispatch/extraction; install.sh's own cron chain pins your login
  shell's claude dir the same way.)

  This launchd plist is provided as a fallback ONLY; it will silently see an empty
  vault if TCC blocks ~/Documents. actd (which never touches ~/Documents) is the
  one that is safe under launchd.
  ==================================================================================

  CREDENTIALS: ANTHROPIC_API_KEY is NOT set here. radar's headless claude reads the
  key from ~/.config/anthropic-key.txt (README landmine #1).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zelin.aiassistant.radar</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>-m</string>
        <string>act.radar</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOURUSERNAME/Projects/zelin-ai-assistant</string>

    <!-- every 30 minutes -->
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AIASSISTANT_HOME</key>
        <string>/Users/YOURUSERNAME/Projects/zelin-ai-assistant</string>
        <key>PATH</key>
        <string>/Users/YOURUSERNAME/.claude-bin:/opt/homebrew/bin:/opt/homebrew/sbin:/Users/YOURUSERNAME/miniconda3/bin:/Users/YOURUSERNAME/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOURUSERNAME/Projects/zelin-ai-assistant/state/radar.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOURUSERNAME/Projects/zelin-ai-assistant/state/radar.launchd.log</string>
</dict>
</plist>
